# files/rootdir/etc/rc.local
#!/bin/sh

iptables -F
iptables -A INPUT -s 127.0.0.1 -p tcp -m tcp --dport 2525 -j ACCEPT
iptables -A INPUT  -p tcp -m tcp --dport 2525 -j DROP

/opt/haraka-inbound/stop-inbound.sh
/opt/haraka-inbound/start-inbound.sh

## block port 2525 of postfix on other interface as its only for internal
