# files/rootdir/var/spool/cron/crontabs/root
MAILTO=""
*/5 * * * * /usr/local/bin/mailwatch-postfix-relay
