# files/rootdir/var/spool/cron/crontabs/root
MAILTO=""
* * * * * /usr/share/webmin/powermail/cron-for-group-mail-list.sh
