# files/rootdir/var/spool/cron/crontabs/root
MAILTO=""
* * * * * /usr/share/webmin/powermail/cron-for-group-mail-list.sh
1 6 * * * /usr/local/src/disk-space-alert.sh
